# nginx/rhixescans.socket
[Unit]
Description=rhixescans socket

[Socket]
ListenStream=/run/rhixescans.sock

[Install]
WantedBy=sockets.target
